# templates/confd.sysv.erb
#!/bin/bash
#
#       /etc/rc.d/init.d/confd
#
#       Daemonize the confd agent.
#
# chkconfig:   2345 95 95
# description: Cnfiguration made easy. \o
#  Distributed, highly available, and datacenter-aware.
# processname: confd
# pidfile: /var/run/confd/pidfile

# Source function library.
. /etc/init.d/functions

CONFD=<%= scope.lookupvar('confd::installdir') %>/confd
PID_FILE=/var/run/confd/confd.pid
LOG_FILE=/var/log/confd.log

[ -e /etc/sysconfig/confd ] && . /etc/sysconfig/confd

export GOMAXPROCS=${GOMAXPROCS:-2}

#
# Create the /var/run/confd directory, which can live on a tmpfs
# filesystem and be destroyed between reboots.
#
mkrundir() {
    [ ! -d /var/run/confd ] && mkdir -p /var/run/confd
}


#
# Create a PID file if it doesn't already exist, for clean upgrades
# from previous init-script controlled daemons.
#
KILLPROC_OPT="-p ${PID_FILE}"
mkpidfile() {
        # Create PID file if it didn't exist
        mkrundir
        [ ! -f $PID_FILE ] && pidofproc $CONFD > $PID_FILE
        if [ $? -ne 0 ] ; then
            rm $PID_FILE
            KILLPROC_OPT=""
        fi
}

start() {
        echo -n "Starting confd: "
        mkrundir
        [ -f $PID_FILE ] && rm $PID_FILE
        daemon --user=root \
            --pidfile="$PID_FILE" \
            "$CONFD" >> "$LOG_FILE" &
        retcode=$?
        success
        echo
        touch /var/lock/subsys/confd
        return $retcode
}

stop() {
        echo -n "Shutting down confd: "
        # If confd is not acting as a server, exit gracefully
        if ("${CONSUL}" info 2>/dev/null | grep -q 'server = false' 2>/dev/null) ; then
            "$CONSUL" leave
        fi

        # If acting as a server, or if leave failed, kill it.
        mkpidfile
        killproc $KILLPROC_OPT $CONSUL -9
        
        echo

        retcode=$?
        rm -f /var/lock/subsys/confd $PID_FILE
        return $retcode
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status confd
        ;;
    restart)
        stop
        start
        ;;
    reload)
        mkpidfile
        killproc $KILLPROC_OPT $CONFD -HUP
        ;;
    condrestart)
        [ -f /var/lock/subsys/confd ] && restart || :
        ;;
    *)
        echo "Usage: confd {start|stop|status|reload|restart}"
        exit 1
        ;;
esac
exit $?
