# templates/systemd.service.erb
[Unit]
Description=ConfD Service
After=network.target

[Service]
ExecStart=<%= scope.lookupvar('::confd::installdir') %>confd
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID

[Install]
WantedBy=multi-user.target
